# pi.dev/linkedin-feed/linkedin-scrape.service
[Unit]
Description=LinkedIn Job Scraper
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/pi.dev/linkedin-feed
ExecStart=/usr/bin/xvfb-run /usr/bin/python3 /opt/pi.dev/linkedin-feed/scrape.py --output /opt/pi.dev/linkedin-feed/jobs_cache.json
StandardOutput=append:/opt/pi.dev/linkedin-feed/scrape.log
StandardError=append:/opt/pi.dev/linkedin-feed/scrape.log
TimeoutStartSec=300
